# nas-enter.service
echo [Unit]
Description=My amazing startup script

[Service]
Type=simple
ExecStart=/usr/bin/nas_enter.sh

[Install]
WantedBy=multi-user.target
